# templates/wildfly.systemd.service.epp
# THIS FILE IS MANAGED BY PUPPET
[Unit]
Description=WildFly application server
Wants=network-online.target
After=network-online.target

[Service]
User=<%= $wildfly::user %>
Group=<%= $wildfly::group %>
EnvironmentFile=/etc/wildfly/wildfly.conf
LimitNOFILE=102642
ExecStart=<%= $wildfly::dirname %>/bin/launch.sh ${WILDFLY_MODE} ${WILDFLY_CONFIG}
Restart=always
RestartSec=20
StandardOutput=null

[Install]
WantedBy=multi-user.target
