# templates/wildfly.sysvinit.service.epp
[Unit]
Description=Wildfly Application Server
After=syslog.target network.target

[Service]
Type=forking
ExecStart=/etc/init.d/<%= $wildfly::service::service_name %> start
ExecStop=/etc/init.d/<%= $wildfly::service::service_name %> stop
<%- if $wildfly::working_dir {-%>
WorkingDirectory=<%= $wildfly::working_dir %>
<%- } -%>

[Install]
WantedBy=multi-user.target
